# install/serial2mqtt.conf
description "Serial Interface to Mosquitto MQTT broker"
author "Joe George <joeg@engimusing.com"

start on net-device-up

respawn

exec /usr/local/bin/serial2mqtt
/etc/init/serial2mqtt.conf (END)
